# configure/supervisor/jybox.net.conf
[program:jybox.net]
command = coffee /home/nodeapp/jybox.net/app.coffee
autorestart = true
directory = /home/nodeapp/jybox.net
user=nodeapp
