# configure/supervisor/newsbee.conf
[program:commitcal]
user = newsbee
command = coffee app.coffee
directory = /home/newsbee/github-commit-ical
autorestart = true
autostart = true
redirect_stderr = true
