# configure/supervisor/newsbee.conf
[program:commitcal]
user = jybox
command = coffee app.coffee
directory = /home/jybox/github-commit-ical
autorestart = true
autostart = true
redirect_stderr = true
